# systemd/dunst-low-battery.timer
[Unit]
Description=Timer to check battery charge

[Timer]
OnActiveSec=8min
OnUnitActiveSec=8min

[Install]
WantedBy=timers.target
