# systemd/dunst-low-battery.service
[Unit]
Description=Check battery charge

[Service]
ExecStart=/bin/sh -c "~/.local/bin/dunst-low-battery"
